# services/com.photoserver.moondream.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.photoserver.moondream</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Library/Frameworks/Python.framework/Versions/3.13/bin/python3</string>
        <string>/Users/thortle/Desktop/media_sorter/services/moondream_service.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/thortle/Desktop/media_sorter</string>

    <key>StandardOutPath</key>
    <string>/Users/thortle/Desktop/media_sorter/services/moondream.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/thortle/Desktop/media_sorter/services/moondream.error.log</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ProcessType</key>
    <string>Background</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PHOTO_DIR</key>
        <string>/Volumes/T7_SSD/G-photos</string>
    </dict>
</dict>
</plist>
